# systemd/allpfit-etl.service
[Unit]
Description=AllpFit Analytics ETL Pipeline V3
After=network.target postgresql.service

[Service]
Type=oneshot
User=isaac
Group=isaac
WorkingDirectory=/home/isaac/projects/allpfit-analytics

# Variáveis de ambiente
Environment="PATH=/home/isaac/projects/allpfit-analytics/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/isaac/projects/allpfit-analytics/.env

# Comando de execução
ExecStart=/home/isaac/projects/allpfit-analytics/venv/bin/python3 \
    /home/isaac/projects/allpfit-analytics/src/features/etl_pipeline_v3.py \
    --triggered-by scheduler

# Restart em caso de falha (max 3 tentativas)
Restart=on-failure
RestartSec=300
StartLimitBurst=3
StartLimitIntervalSec=3600

# Logging
StandardOutput=append:/home/isaac/projects/allpfit-analytics/logs/etl/systemd.log
StandardError=append:/home/isaac/projects/allpfit-analytics/logs/etl/systemd_error.log

# Timeout
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
