# systemd/allpfit-etl.timer
[Unit]
Description=Timer para AllpFit Analytics ETL (executa diariamente às 3h)
Requires=allpfit-etl.service

[Timer]
# Executar diariamente às 3:00 AM
OnCalendar=*-*-* 03:00:00

# Aguardar 5 minutos se o sistema estiver desligado na hora agendada
Persistent=true

# Se perdeu a execução, executar assim que possível
AccuracySec=1min

[Install]
WantedBy=timers.target
