# systemd/etl-allpfit.timer
[Unit]
Description=ETL AllpFit Analytics - Execução Automática a cada 2 horas
Documentation=https://github.com/allpfit/analytics
Requires=etl-allpfit.service

[Timer]
# Executar a cada 2 horas
OnCalendar=0/2:00:00

# Se o sistema estava desligado no horário agendado, executar assim que ligar
Persistent=true

# Adicionar aleatoriedade de até 5 minutos (evitar sobrecarga)
RandomizedDelaySec=5min

# Unidade a ser iniciada
Unit=etl-allpfit.service

[Install]
WantedBy=timers.target